# local/systemd/homemonitor.service
[Unit]
Description=homemonitor
Wants=homemonitor.timer

[Service]
Type=simple
ExecStart=/usr/local/bin/homemonitor
