# local/systemd/pingcheck.service
[Unit]
Description=pingcheck
Wants=pingcheck.timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pingcheck
