# local/systemd/pingcheck.timer
[Unit]
Description=Runs pingcheck every hour

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
Unit=pingcheck.service

[Install]
WantedBy=multi-user.target
