# local/systemd/timelapse.service
[Unit]
Description=TimeLapse
Wants=timelapse.timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/make_timelapse
